# .config/systemd/user/joplin.timer
[Unit]
Description=Joplin resource timer

[Timer]
OnBootSec=0min
OnCalendar=*:*:*
Persistent=true

[Install]
WantedBy=timers.target
